# provision/com.ledatic.ledaticground-iq.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ledatic.ledaticground-iq</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/ledaticempire/projects/ledaticground/scripts/pass_scheduler.sh</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>CAPTURE</key>
        <string>iq</string>
        <key>MINEL</key>
        <string>40</string>
        <key>GAIN</key>
        <string>49</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>60</integer>
    <key>StandardOutPath</key>
    <string>/tmp/iqsched.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/iqsched.log</string>
</dict>
</plist>
